# dir_root/etc/profile.d/20_rosterlib.sh
#!/usr/bin/env bash
##
#
## Roster lib, here you call all the libraries to load in the session

. ${BASHTUI_LIB}/say.sh
. ${BASHTUI_LIB}/svnlocal.sh
. ${BASHTUI_LIB}/custom_functions.sh
